# ops/systemd/sotsiaalai-research-worker.service
[Unit]
Description=SotsiaalAI Research Worker
After=network-online.target sotsiaalai-rag.service
Wants=network-online.target sotsiaalai-rag.service

[Service]
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/apps/sotsiaalai
Environment=NODE_ENV=production
Environment=PATH=/usr/bin:/bin
Environment=RESEARCH_WORKER_ID=research-worker-1
EnvironmentFile=/etc/sotsiaalai/frontend.env
ExecStart=/usr/bin/node /home/ubuntu/apps/sotsiaalai/node_modules/tsx/dist/cli.mjs --tsconfig jsconfig.json scripts/research-worker.mjs
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
