# ops/livekit/livekit.service
[Unit]
Description=SotsiaalAI self-hosted LiveKit server
Documentation=https://docs.livekit.io/transport/self-hosting/deployment/
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=livekit
Group=livekit
ExecStart=/usr/local/bin/livekit-server --config /etc/livekit/livekit.yaml
Restart=always
RestartSec=5
LimitNOFILE=1048576
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
